# smf/manifests/buckets-mdapi.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright 2020 Joyent, Inc.
-->

<service_bundle type="manifest" name="buckets-mdapi">
    <service name="manta/application/buckets-mdapi" type="service" version="1">
        <create_default_instance enabled="true" />

        <dependency name="network"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
            <service_fmri value="svc:/network/physical" />
        </dependency>

        <dependency name="filesystem"
                    grouping="require_all"
                    restart_on="error"
                    type="service">
            <service_fmri value="svc:/system/filesystem/local" />
        </dependency>

        <dependency name="mdata"
                    grouping="require_all"
                    restart_on="none"
                    type="service">
            <service_fmri value="svc:/smartdc/mdata:execute" />
        </dependency>

        <dependency name="buckets-mdapi-setup"
                    grouping="optional_all"
                    restart_on="none"
                    type="service">
            <service_fmri value="svc:/manta/application/buckets-mdapi-setup" />
        </dependency>

        <exec_method type="method" name="start" timeout_seconds="10"
          exec="./bin/buckets-mdapi -c ./etc/config.toml &amp;">
          <method_context working_directory="/opt/smartdc/buckets-mdapi">
                  <method_credential user="nobody"
                                   group="nobody"
                                   privileges="basic,net_privaddr"/>
            </method_context>
        </exec_method>

        <exec_method type="method" name="stop" exec=":kill"
          timeout_seconds="30" />

        <template>
            <common_name>
                <loctext xml:lang="C">Manta Buckets Metadata API Service</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
